# extras/linux/systemd/nabud.service.in
[Unit]
Description=A server for the NABU PC
Wants=network-online.target
After=network-online.target

[Service]
Type=idle
ExecStart=XXXsbindirXXX/nabud -u nabu -f -l /var/log/nabud
WorkingDirectory=/

[Install]
WantedBy=multi-user.target
